# com.behilnortu.pdf-manager.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.behilnortu.pdf-manager</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/Hubie1/.nvm/versions/node/v22.22.0/bin/node</string>
        <string>/Users/Hubie1/Desktop/Claude Code/server/server.js</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/Hubie1/Desktop/Claude Code</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/Hubie1/Library/Logs/pdf-manager.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/Hubie1/Library/Logs/pdf-manager-error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/Hubie1/.nvm/versions/node/v22.22.0/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>NODE_ENV</key>
        <string>production</string>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
